# profile.d/variables.sh
export RAILS_ENV="production"
export LC_ALL="en_US.UTF-8"
